# profile.d/environment.sh
export PREFIX=/data/data/com.termux/files/usr
export SVDIR=/data/adb/chroot-service/containers
export LOGDIR=/cache
export PATH=$PREFIX/bin:$PATH
