# etc/systemd/system/android-network-monitor.service
[Unit]
Description=Android Network Monitor Service

[Service]
Type=simple
ExecStart=/usr/local/bin/android-network-monitor
Restart=always

[Install]
WantedBy=multi-user.target